# jdeb/src/main/deb/control/prerm
#!/bin/sh

printf "Running prerm script\n"

update-rc.d -f jdeb-tips remove

